# launchagents/com.walter.weekly-cache-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.walter.weekly-cache-cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/walter/.local/bin/weekly-cache-cleanup.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/walter/Library/Logs/weekly-cache-cleanup.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/walter/Library/Logs/weekly-cache-cleanup.launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
